# conf/supervisor/server-side-render.conf
[group:talebook]
programs=tornado,nginx,nodejs

[program:tornado]
command=gosu talebook:talebook python3 server.py --with-library=/data/books/library --port=8000 --host=127.0.0.1 --logging=debug --log-file-prefix=/data/log/talebook.log
directory=/var/www/talebook
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:nginx]
command = /usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=unexpected

[program:nodejs]
command=gosu talebook:talebook npm run start
directory=/var/www/talebook/app
autorestart=true
redirect_stderr=true
stdout_logfile=/data/log/nodejs.log

