# conf/supervisor/talebook.conf
[group:talebook]
programs=tornado-8000

[program:tornado-8000]
command=python3 talebook/server.py --with-library=/data/books/library --port=8000 --host=127.0.0.1 --logging=debug
directory=/var/www
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/data/log/talebook.log
loglevel=info

